# com.maedashiori.design-md-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.maedashiori.design-md-update</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/maedashiori/awesome-design-md/run_auto_update.sh</string>
    </array>

    <!-- 毎月1日 9:00 に実行 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/maedashiori/awesome-design-md/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/maedashiori/awesome-design-md/launchd.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
